# scripts/supervisord.conf
[supervisord]
nodaemon=true
logfile=/app/logs/supervisord.log
pidfile=/tmp/supervisord.pid

[program:smtp]
command=python /app/scripts/smtp_server.py
autostart=true
autorestart=true
stdout_logfile=/app/logs/smtp.log
stderr_logfile=/app/logs/smtp.log
environment=SIDECAR_HOST="%(ENV_SIDECAR_HOST)s",SMTP_HOST="%(ENV_SMTP_HOST)s",SMTP_PORT="%(ENV_SMTP_PORT)s",SMTP_USER="%(ENV_SMTP_USER)s",SMTP_PASS="%(ENV_SMTP_PASS)s",SMTP_FROM="%(ENV_SMTP_FROM)s"

[program:imap]
command=python /app/scripts/imap_server.py
autostart=true
autorestart=true
stdout_logfile=/app/logs/imap.log
stderr_logfile=/app/logs/imap.log
environment=IMAP_HOST="%(ENV_IMAP_HOST)s",IMAP_PORT="%(ENV_IMAP_PORT)s",IMAP_USER="%(ENV_IMAP_USER)s",IMAP_PASS="%(ENV_IMAP_PASS)s",IMAP_FOLDER="%(ENV_IMAP_FOLDER)s",IMAP_POLL_SEC="%(ENV_IMAP_POLL_SEC)s",N8N_INBOUND_URL="%(ENV_N8N_INBOUND_URL)s"

[program:recording]
command=python /app/scripts/recording_server.py
autostart=true
autorestart=true
stdout_logfile=/app/logs/recording.log
stderr_logfile=/app/logs/recording.log
environment=SIDECAR_HOST="%(ENV_SIDECAR_HOST)s",RECORDINGS_DIR="/app/recordings"

[program:auth]
command=python /app/scripts/auth_server.py
autostart=true
autorestart=true
stdout_logfile=/app/logs/auth.log
stderr_logfile=/app/logs/auth.log
environment=SIDECAR_HOST="%(ENV_SIDECAR_HOST)s",PG_HOST="%(ENV_PG_HOST)s",PG_PORT="%(ENV_PG_PORT)s",PG_USER="%(ENV_PG_USER)s",PG_PASSWORD="%(ENV_PG_PASSWORD)s",PG_DB="%(ENV_PG_DB)s",SEED_ADMIN_EMAIL="%(ENV_SEED_ADMIN_EMAIL)s",SEED_ADMIN_PASSWORD="%(ENV_SEED_ADMIN_PASSWORD)s"
